# prototype.conf
#
# This is the upstart config for the workflow prototype application.
#

start on runlevel [2345]
stop on runlevel [016]

# NB: setuid is only supported on recent versions of upstart - i.e.
#  on ubuntu not on amazon linux

setuid workflow

chdir /home/workflow

# automatically restart if the process dies
# respawn

env USER_HOME=/home/workflow

env APP=prototype
env JAVA_OPTS="-Dfile.encoding=UTF-8 -Dconfig.resource=application.@STAGE@.conf -Xmx300m -XX:MaxPermSize=64m -XX:+PrintGCDetails -XX:+PrintGCDateStamps -Xloggc:gc.log "

env AWS_ACCESS_KEY=@AWS_ACCESS_KEY@
env AWS_SECRET=@AWS_SECRET@

env AWS_FLEX_NOTIFICATIONS_QUEUE=@AWS_FLEX_NOTIFICATIONS_QUEUE@

env COMPOSER_URL=@COMPOSER_URL@
env PRESENCE_URL=@PRESENCE_URL@

env DB_PORT=@DB_PORT@
env DB_URL=@DB_URL@
env DB_PASSWORD=@DB_PASSWORD@
env DB_USER=@DB_USER@
env GOOGLE_CLIENT_ID=@GOOGLE_CLIENT_ID@
env GOOGLE_SECRET=@GOOGLE_SECRET@
env HOST=@HOST@
env PANDOMAIN_DOMAIN=@PANDOMAIN_DOMAIN@
env PANDOMAIN_AWS_KEY=@PANDOMAIN_AWS_KEY@
env PANDOMAIN_AWS_SECRET=@PANDOMAIN_AWS_SECRET@
env PRESENCE_CLIENTLIB=@PRESENCE_CLIENTLIB@
env LOGFILE=/home/workflow/logs/stdout.log

script
  PROTOTYPE_CMD="$USER_HOME/$APP/bin/$APP"
  echo "$PROTOTYPE_CMD" > $USER_HOME/cmd.txt
  $PROTOTYPE_CMD > $LOGFILE 2>&1
end script
